# s6/debian-root/etc/services.d/php7.3-fpm/run
#!/usr/bin/execlineb -P
service php7.3-fpm start